# systemd/example.target
[Unit]
Description=Example BEAM application with socket activation
Requires=example-activation.socket
