# systemd/example-epmd.socket
[Unit]
Requires=sockets.target

[Socket]
ListenStream=4369
BindIPv6Only=both
ReusePort=true
