# systemd/example-activation.service
[Unit]
Description=Proxy service
Requires=example-application.socket

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd /var/example-service
